# supervisord.conf
[supervisord]
nodaemon=true

[program:fixmystreet]
command=/etc/init.d/fixmystreet start

[program:nginx]
command=nginx